# scripts/systemd/softtarget.service
[Unit]
Description=Soft Target Backend (FastAPI + Gunicorn)
After=network.target postgresql.service
Wants=postgresql.service

[Service]
Type=simple
User=softtarget
Group=softtarget
WorkingDirectory=/opt/softtarget
EnvironmentFile=/etc/softtarget/softtarget.env
ExecStart=/opt/softtarget/.venv/bin/gunicorn app.main:app --config /opt/softtarget/scripts/gunicorn.conf.py
ExecReload=/bin/kill -s HUP $MAINPID
Restart=on-failure
RestartSec=5s
KillSignal=SIGTERM
TimeoutStopSec=30

# Hardening
NoNewPrivileges=true
PrivateTmp=true
PrivateDevices=true
ProtectSystem=strict
ProtectHome=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
RestrictSUIDSGID=true
RestrictNamespaces=true
LockPersonality=true
ReadWritePaths=/var/lib/softtarget /var/log/softtarget

[Install]
WantedBy=multi-user.target
